# fs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh
/usr/sbin/crond -f -d 8 -l 9 -L /dev/stdout
