# fs/etc/crontabs/root
# Reload mbtiles (if required) every minute.
* * * * * /usr/local/bin/reload-mbtiles-if-changed 1>/dev/stdout 2>/dev/stderr
