# fs/etc/services.d/mbtileserver/run
#!/usr/bin/with-contenv sh
exec /mbtileserver --enable-reload-signal $MBTILESERVER_ARGS
