# fs/etc/services.d/mbtileserver/run
#!/usr/bin/with-contenv sh
exec /mbtileserver --disable-preview --disable-svc-list --enable-reload-signal --root-url /tilesets
